# debian/tron.service
[Unit]
Description=trond
After=network.target

[Service]
User=tron
EnvironmentFile=/etc/default/tron
ExecStart=/usr/bin/trond --lock-file=${LOCKFILE:-$PIDFILE} --working-dir=${WORKINGDIR} --host ${LISTEN_HOST} --port ${LISTEN_PORT} ${DAEMON_OPTS}
TimeoutStopSec=20
Restart=on-failure
LimitNOFILE=10000
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
